# config/com.local.mlx-server-restart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.mlx-server-restart</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/launchctl</string>
        <string>kickstart</string>
        <string>-k</string>
        <string>gui/501/com.local.mlx-server</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/mlx-restart.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/mlx-restart.log</string>
</dict>
</plist>
